# systemd/.config/systemd/user/egyetem-backup.timer
[Unit]
Description=Backup timer

[Timer]
Unit=egyetem-backup.service
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
